# init/tourboxneo.service
[Unit]
Description=TourBox NEO

[Service]
ExecStart=/usr/bin/tourboxneo

[Install]
WantedBy=default.target
